# deploy/parking-app.service
# Systemd service file for Parking App Django backend
# Place this file at: /etc/systemd/system/parking-app.service

[Unit]
Description=Parking in a Pinch Django App
Requires=parking-app.socket
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=parking-app
WorkingDirectory=/var/www/parking-app/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings.production
EnvironmentFile=/var/www/parking-app/backend/.env
ExecStart=/var/www/parking-app/backend/venv/bin/gunicorn \
    --config /var/www/parking-app/deploy/gunicorn.conf.py \
    config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStopSec=5

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/parking-app/backend/media /var/log/parking-app /var/run/parking-app

[Install]
WantedBy=multi-user.target